# deploy/supervisor_gunicorn.conf
[program:gunicorn]
command=/home/ubuntu/tenants_back/bin/gunicorn_start.sh
; No user= here on purpose: supervisor launches the script as root, and
; gunicorn itself drops privileges via --user=ubuntu --group=www-data. That
; lets it set the process group to www-data (a non-root start could not, since
; supplementary group membership is not enough for setgid) and chown the socket
; to ubuntu:www-data so nginx can connect.
autostart=true
autorestart=true
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/var/log/gunicorn.log
stderr_logfile=/var/log/gunicorn.err.log
environment=PYTHONUNBUFFERED=1
