# deploy/supervisor_celery.conf
; Celery for the MVP single host. Three workers (concurrency 2 each) + one beat:
;   celery_fast          -> only the `fast` queue
;   celery_slow          -> only the `slow` queue
;   celery_slow_service  -> `slow` + `service` queues
;   celery_beat          -> the scheduler (publishes due periodic tasks)
;
; Notes:
; * `%%h` is escaped here: supervisor treats a lone `%` as its own expansion
;   syntax, so Celery's hostname placeholder `%h` must be written `%%h`.
; * Each worker needs a UNIQUE node name (-n ...@%%h) when several run on one host.
; * Logs: supervisor captures each process's stdout+stderr into its own file under
;   the project's logs/ dir (ensure it exists: `mkdir -p logs`).
; * stopasgroup/killasgroup: a prefork worker spawns child processes; signal the
;   whole group so children stop too. stopwaitsecs gives running tasks time to finish.
; * celery_beat must run as EXACTLY ONE instance across the whole deployment. On a
;   single host that's this one program; do NOT also start beat elsewhere (two beats
;   => duplicate per-tenant fan-out). It uses the DB scheduler, so no --schedule file.

[program:celery_fast]
command=/home/ubuntu/mt-shards-back-test/venv/bin/celery -A tenants_back worker -Q fast --concurrency=2 -n fast@%%h --loglevel=info
directory=/home/ubuntu/mt-shards-back-test
user=ubuntu
numprocs=1
stdout_logfile=/home/ubuntu/mt-shards-back-test/logs/celery_fast.log
redirect_stderr=true
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=120
stopasgroup=true
killasgroup=true
environment=PYTHONUNBUFFERED=1,LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8

[program:celery_slow]
command=/home/ubuntu/mt-shards-back-test/venv/bin/celery -A tenants_back worker -Q slow --concurrency=2 -n slow@%%h --loglevel=info
directory=/home/ubuntu/mt-shards-back-test
user=ubuntu
numprocs=1
stdout_logfile=/home/ubuntu/mt-shards-back-test/logs/celery_slow.log
redirect_stderr=true
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=300
stopasgroup=true
killasgroup=true
environment=PYTHONUNBUFFERED=1,LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8

[program:celery_slow_service]
command=/home/ubuntu/mt-shards-back-test/venv/bin/celery -A tenants_back worker -Q slow,service --concurrency=2 -n slowsvc@%%h --loglevel=info
directory=/home/ubuntu/mt-shards-back-test
user=ubuntu
numprocs=1
stdout_logfile=/home/ubuntu/mt-shards-back-test/logs/celery_slow_service.log
redirect_stderr=true
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=300
stopasgroup=true
killasgroup=true
environment=PYTHONUNBUFFERED=1,LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8

[program:celery_beat]
command=/home/ubuntu/mt-shards-back-test/venv/bin/celery -A tenants_back beat --loglevel=info
directory=/home/ubuntu/mt-shards-back-test
user=ubuntu
numprocs=1
stdout_logfile=/home/ubuntu/mt-shards-back-test/logs/celery_beat.log
redirect_stderr=true
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
stopasgroup=true
killasgroup=true
environment=PYTHONUNBUFFERED=1,LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8

[group:celery]
programs=celery_fast,celery_slow,celery_slow_service,celery_beat
